# src/librootfs/akaimer39lib/etc/cron/crontabs/root
00 00 * * * if pidof camera.sh; then reboot; fi
