# src/librootfs/akaimer39lib/etc/init.d/rc.local
#!/bin/sh

# Mount sd card
mount /dev/mmcblk0p1 /mnt

# Adding Hwclock to the system
echo ds1307 0x68 > /sys/class/i2c-adapter/i2c-0/new_device
hwclock --hctosys

if [ ! -f /etc/mtab ]; then ln -s /proc/mounts /etc/mtab; fi

echo "5120" > /proc/sys/vm/min_free_kbytes

chmod 600 .rsync

ifplugd -an -t 5 -u1 -d4 -i usb0 &
ifup usb0

crond

if [ -f /mnt/interfaces ]
then
/etc/init.d/setup.sh &
fi
